# .config/systemd/user/hyprlock.service
[Unit]
Description=Hyprlock screen locker
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/hyprlock
